# packaging/inferno/prerm
#!/bin/sh
set -e

if [ "${1:-}" = "remove" ] || [ "${1:-}" = "deconfigure" ]; then
  systemctl disable --now teletool-inferno-clock.service 2>/dev/null || true
fi

exit 0
